# resources/mpdris.service
[Unit]
Description=Music Player Daemon MPRIS bridge
BindsTo=mpd.service
After=mpd.service

[Service]
Type=notify-reload
NotifyAccess=main
Restart=on-failure
ExecStart=/usr/local/bin/mpDris --service

[Install]
RequiredBy=mpd.service
